# files/etc/profile.d/anaconda.sh
ANACONDA_DIR=/opt/anaconda
PATH=${ANACONDA_DIR}/bin:${PATH}
PYTHON=${ANACONDA_DIR}/bin/python
LD_LIBRARY_PATH=${ANACONDA_DIR}/lib:${LD_LIBRARY_PATH}

