# files/etc/profile.d/anaconda.sh
ANACONDA_DIR=/opt/anaconda
PATH=${ANACONDA_DIR}/bin:${PATH}
PYTHON=${ANACONDA_DIR}/bin/python
