# files/etc/profile.d/anaconda.sh
export ANACONDA_DIR=/opt/anaconda
export PATH=${ANACONDA_DIR}/bin:${PATH}
export PYTHON=${ANACONDA_DIR}/bin/python
